# backend-implementation/docker/siem/osquery/config/osquery.conf
# ============================================================================
# OSQUERY HOST-BASED MONITORING CONFIGURATION
# ============================================================================
#
# Production OSQuery configuration for endpoint security monitoring
# SQL-based queries for system state and security analysis
#
# Created by: DEVOPS-AGENT - TIER 1 Advanced Threat Protection
# Date: 2025-08-14
# Version: 1.0.0
# ============================================================================

--config_plugin=filesystem
--logger_plugin=filesystem
--database_path=/var/log/osquery/osquery.db
--results_log_file=/var/log/osquery/osqueryd.results.log
--snapshot_log_file=/var/log/osquery/osqueryd.snapshots.log
--health_log_file=/var/log/osquery/osqueryd.health.log
--log_result_events=true
--verbose=false
--worker_threads=4
--enable_monitor=true
--config_refresh=300
--schedule_splay_percent=10
--events_expiry=3600
--purge_inbox=true
--inbox_max=100000
--read_max=50000000
--events_optimize=true
--events_max=50000
--disable_events=false
--audit_allow_config=true
--audit_allow_sockets=true
--enable_syslog=false
--host_identifier=hostname
--schedule_default_interval=3600

{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem",
    "logger_path": "/var/log/osquery",
    "database_path": "/var/log/osquery/osquery.db",
    "log_result_events": true,
    "schedule_splay_percent": 10,
    "pidfile": "/var/osquery/osquery.pidfile",
    "events_expiry": 3600,
    "database_dump": true,
    "verbose": false,
    "worker_threads": 4,
    "enable_monitor": true,
    "disable_distributed": false,
    "distributed_interval": 60,
    "distributed_plugin": "tls",
    "distributed_get_queries_retries": 3,
    "distributed_write_results_retries": 3,
    "host_identifier": "hostname",
    "schedule_default_interval": 3600,
    "audit_allow_config": true,
    "audit_allow_sockets": true,
    "audit_allow_process_events": true,
    "audit_allow_fim_events": true,
    "enable_keyboard_events": false,
    "enable_mouse_events": false,
    "disable_events": false,
    "events_optimize": true,
    "events_max": 50000,
    "read_max": 50000000,
    "carver_disable_function": false,
    "carver_start_endpoint": "",
    "carver_continue_endpoint": "",
    "carver_block_size": 2000000
  },

  "schedule": {
    "system_info": {
      "query": "SELECT hostname, cpu_brand, hardware_model, hardware_vendor, hardware_serial, physical_memory FROM system_info;",
      "interval": 3600,
      "description": "Basic system information"
    },

    "osquery_info": {
      "query": "SELECT pid, version, config_hash, config_valid, extensions, build_platform, build_distro FROM osquery_info;",
      "interval": 3600,
      "description": "OSQuery daemon information"
    },

    "kernel_info": {
      "query": "SELECT version, arguments FROM kernel_info;",
      "interval": 3600,
      "description": "Kernel and boot arguments"
    },

    "os_version": {
      "query": "SELECT name, version, major, minor, patch, build, platform, platform_like, arch FROM os_version;",
      "interval": 3600,
      "description": "Operating system version information"
    },

    "uptime": {
      "query": "SELECT days, hours, minutes, seconds, total_seconds FROM uptime;",
      "interval": 300,
      "description": "System uptime information"
    },

    "processes": {
      "query": "SELECT pid, name, path, cmdline, cwd, root, uid, gid, euid, egid, on_disk, wired_size, resident_size, total_size, user_time, system_time, disk_bytes_read, disk_bytes_written, start_time, parent, pgroup, threads, nice, elapsed_time FROM processes;",
      "interval": 300,
      "description": "Current running processes"
    },

    "process_open_sockets": {
      "query": "SELECT DISTINCT process.name, listening.port, listening.address, process.pid FROM processes AS process JOIN listening_ports AS listening ON process.pid = listening.pid WHERE listening.address != '127.0.0.1' AND listening.address != '::1';",
      "interval": 300,
      "description": "Processes with open network sockets"
    },

    "network_interfaces": {
      "query": "SELECT interface, address, mask, broadcast, point_to_point, type, mtu, metric, flags, ipackets, opackets, ibytes, obytes, ierrors, oerrors, idrops, odrops, collisions, last_change FROM interface_details;",
      "interval": 600,
      "description": "Network interface information"
    },

    "listening_ports": {
      "query": "SELECT pid, port, protocol, family, address FROM listening_ports;",
      "interval": 300,
      "description": "Network listening ports"
    },

    "arp_cache": {
      "query": "SELECT address, mac, interface, permanent FROM arp_cache;",
      "interval": 300,
      "description": "ARP cache entries"
    },

    "route_table": {
      "query": "SELECT destination, netmask, gateway, source, flags, interface, mtu, metric, type FROM routes;",
      "interval": 600,
      "description": "Network routing table"
    },

    "installed_software": {
      "query": "SELECT name, version, source, pid FROM rpm_packages UNION SELECT name, version, source, '' AS pid FROM deb_packages;",
      "interval": 3600,
      "description": "Installed software packages"
    },

    "kernel_modules": {
      "query": "SELECT name, size, used_by, status, address FROM kernel_modules;",
      "interval": 600,
      "description": "Loaded kernel modules"
    },

    "startup_items": {
      "query": "SELECT name, path, args, type, source, status, username FROM startup_items;",
      "interval": 3600,
      "description": "System startup items"
    },

    "crontab": {
      "query": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
      "interval": 3600,
      "description": "Cron job schedules"
    },

    "users": {
      "query": "SELECT uid, gid, username, description, directory, shell, uuid FROM users;",
      "interval": 3600,
      "description": "Local user accounts"
    },

    "groups": {
      "query": "SELECT gid, groupname, comment FROM groups;",
      "interval": 3600,
      "description": "Local groups"
    },

    "user_groups": {
      "query": "SELECT uid, gid FROM user_groups;",
      "interval": 3600,
      "description": "User group memberships"
    },

    "logged_in_users": {
      "query": "SELECT user, tty, host, time, pid FROM logged_in_users;",
      "interval": 300,
      "description": "Currently logged in users"
    },

    "last": {
      "query": "SELECT username, tty, pid, type, type_name, time, host FROM last LIMIT 50;",
      "interval": 300,
      "description": "Recent user logins"
    },

    "ssh_config": {
      "query": "SELECT block, host, hostname, port, user, protocol, cipher, compression, timeout, permit_root_login, permit_empty_passwords, password_authentication, pubkey_authentication, challenge_response_authentication, use_pam, x11_forwarding, use_dns FROM ssh_config;",
      "interval": 3600,
      "description": "SSH configuration"
    },

    "file_hashes": {
      "query": "SELECT path, md5, sha1, sha256 FROM hash WHERE path IN ('/bin/sh', '/bin/bash', '/usr/bin/sudo', '/etc/passwd', '/etc/shadow', '/etc/ssh/sshd_config', '/etc/sudoers');",
      "interval": 3600,
      "description": "Critical file hashes"
    },

    "docker_containers": {
      "query": "SELECT id, name, image, image_id, command, created, state, status FROM docker_containers;",
      "interval": 300,
      "description": "Docker container information"
    },

    "docker_images": {
      "query": "SELECT id, created, size_bytes, tags FROM docker_images;",
      "interval": 600,
      "description": "Docker image information"
    },

    "docker_networks": {
      "query": "SELECT id, name, driver, created, enable_ipv6, ipam_driver, ipam_subnet, ipam_gateway FROM docker_networks;",
      "interval": 600,
      "description": "Docker network configuration"
    },

    "docker_volumes": {
      "query": "SELECT name, driver, mountpoint, created FROM docker_volumes;",
      "interval": 600,
      "description": "Docker volume information"
    },

    "memory_info": {
      "query": "SELECT memory_total, memory_free, memory_available, buffers, cached, swap_cached, active, inactive, swap_total, swap_free FROM memory_info;",
      "interval": 300,
      "description": "System memory information"
    },

    "disk_usage": {
      "query": "SELECT device, mount, type, blocks_size, blocks, blocks_free, blocks_available, inodes, inodes_free, flags FROM disk_usage;",
      "interval": 300,
      "description": "Disk usage information"
    },

    "cpu_time": {
      "query": "SELECT core, user, nice, system, idle, iowait, irq, softirq, steal, guest, guest_nice FROM cpu_time;",
      "interval": 300,
      "description": "CPU usage statistics"
    },

    "load_average": {
      "query": "SELECT period, average FROM load_average;",
      "interval": 300,
      "description": "System load average"
    },

    "hardware_events": {
      "query": "SELECT device, model, serial, description, class, driver FROM hardware_events;",
      "interval": 600,
      "description": "Hardware events"
    },

    "usb_devices": {
      "query": "SELECT usb_address, usb_port, vendor, vendor_id, model, model_id, class, subclass, protocol, removable FROM usb_devices;",
      "interval": 600,
      "description": "USB device information"
    },

    "pci_devices": {
      "query": "SELECT pci_slot, driver, vendor, vendor_id, model, model_id, class, subclass, pci_class, pci_subclass FROM pci_devices;",
      "interval": 600,
      "description": "PCI device information"
    },

    "suid_bin": {
      "query": "SELECT path, mode, uid, gid FROM suid_bin;",
      "interval": 3600,
      "description": "SUID/SGID executables"
    },

    "process_memory_map": {
      "query": "SELECT pid, start, end, permissions, offset, device, inode, path, pseudo FROM process_memory_map WHERE pid IN (SELECT pid FROM processes WHERE name IN ('postgres', 'nginx', 'node', 'redis-server', 'elasticsearch'));",
      "interval": 600,
      "description": "Memory map for critical processes"
    },

    "process_file_accesses": {
      "query": "SELECT pid, path, fd FROM process_open_files WHERE pid IN (SELECT pid FROM processes WHERE name IN ('postgres', 'nginx', 'node', 'redis-server', 'elasticsearch'));",
      "interval": 300,
      "description": "File accesses for critical processes"
    },

    "certificates": {
      "query": "SELECT common_name, organization, organization_unit, serial_number, issuer_common_name, issuer_organization, issuer_organization_unit, valid_from, valid_to, fingerprint, key_algorithm, key_strength, key_usage, subject_key_id, authority_key_id, sha1_fingerprint, path FROM certificates WHERE path LIKE '/etc/ssl/certs/%' OR path LIKE '/usr/share/ca-certificates/%';",
      "interval": 3600,
      "description": "SSL/TLS certificates"
    }
  },

  "file_paths": {
    "etc": [
      "/etc/%%"
    ],
    "home": [
      "/home/%%"
    ],
    "tmp": [
      "/tmp/%%"
    ],
    "usr": [
      "/usr/%%"
    ],
    "opt": [
      "/opt/%%"
    ],
    "var": [
      "/var/%%"
    ]
  },

  "exclude_paths": {
    "etc": [
      "/etc/sysconfig/iptables-save",
      "/etc/sysconfig/ip6tables-save"
    ],
    "tmp": [
      "/tmp/osquery%",
      "/tmp/.%"
    ]
  },

  "events": {
    "disable_subscribers": ["user_events"],
    "expiry_time": 3600,
    "optimize": true
  },

  "decorators": {
    "load": [
      "SELECT uuid AS host_uuid FROM system_info;",
      "SELECT user AS username FROM logged_in_users WHERE user <> '' ORDER BY time DESC LIMIT 1;"
    ],
    "always": [
      "SELECT total_seconds AS uptime FROM uptime;"
    ],
    "interval": {
      "3600": [
        "SELECT version FROM kernel_info;"
      ]
    }
  },

  "packs": {
    "incident-response": {
      "discovery": [
        "SELECT pid FROM processes WHERE name = 'osqueryd';"
      ],
      "platform": "posix",
      "version": "1.4.5",
      "queries": {
        "running_processes": {
          "query": "select pid, parent, name, path, cmdline, cwd, root, uid, gid, start_time, user_time, system_time from processes;",
          "interval": 900,
          "description": "Snapshot of running processes"
        },
        "network_connections": {
          "query": "select pid, fd, family, protocol, local_address, local_port, remote_address, remote_port from process_open_sockets where local_port != 0 OR remote_port != 0;",
          "interval": 300,
          "description": "Process network connections"
        }
      }
    },

    "ossec-rootkit": {
      "platform": "linux",
      "version": "2.9.0",
      "discovery": [
        "SELECT name FROM kernel_modules WHERE name LIKE '%rootkit%';"
      ],
      "queries": {
        "rootkit_files": {
          "query": "select * from file where path in ('/tmp/.ICE-unix/1', '/tmp/.X11-unix/1', '/tmp/.font-unix/1', '/dev/.udev', '/usr/share/.config', '/var/run/.lock', '/etc/.pwd.lock', '/etc/.shadow.lock');",
          "interval": 3600,
          "description": "Common rootkit file locations"
        },
        "rootkit_processes": {
          "query": "select * from processes where name in ('xchk', 'chkproc', 'aspyr', 'adore', 'kit', 'hdparm', 'sensord');",
          "interval": 300,
          "description": "Common rootkit process names"
        }
      }
    }
  },

  "yara": {
    "signatures": {
      "/tmp": [
        "/etc/osquery/yara/malware.yar"
      ],
      "/home": [
        "/etc/osquery/yara/malware.yar"
      ]
    }
  }
}